# my-docker-env/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; 默认所有子进程以 root 运行，但下面每个 program 单独降权

[program:openclaw-gateway]
command=/workspace/start-openclaw.sh
directory=/home/ubuntu
; 关键：以 ubuntu 用户运行
user=ubuntu
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopwaitsecs=30
; 环境变量
environment=HOME="/home/ubuntu",NODE_ENV="production",OPENCLAW_CONFIG_PATH="/home/ubuntu/.openclaw/openclaw.json"
; 日志文件由 Supervisor（root）创建，但 openclaw 可写
stdout_logfile=/var/log/supervisor/openclaw-gateway.out.log
stderr_logfile=/var/log/supervisor/openclaw-gateway.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3


;[program:code-server]
;directory=/home/ubuntu
;user=ubuntu
;environment=HOME="/home/ubuntu",PASSWORD="0cccbedfb661cefb0cea5ad3b866b75d1106fdcdccc12fdd",PATH="/usr/local/bin:/usr/bin:/bin:/home/ubuntu/.local/bin"
;command=/workspace/start-code-server.sh
;autorestart=true
;startsecs=5
;startretries=3
;priority=200
;stdout_logfile=/var/log/supervisor/code-server.log
;stderr_logfile=/var/log/supervisor/code-server_err.log
;stdout_logfile_maxbytes=10MB
;stderr_logfile_maxbytes=10MB

[program:opencode]
directory=/home/ubuntu
user=ubuntu
environment=HOME="/home/ubuntu",OPENCODE_SERVER_PASSWORD="0cccbedfb661cefb0cea5ad3b866b75d1106fdcdccc12fdd",PATH="/usr/local/bin:/usr/bin:/bin:/home/ubuntu/.local/bin"
command=/workspace/start-opencode.sh
autorestart=true
startsecs=5
startretries=3
priority=300
stdout_logfile=/var/log/supervisor/opencode.log
stderr_logfile=/var/log/supervisor/opencode_err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB


